# core/services/nexus-health.timer
[Unit]
Description=NexusOS Hourly Health Check Timer
Documentation=https://github.com/nexus-os/nexus-os

[Timer]
OnCalendar=hourly
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
