# core/media-stack/homelab/infrastructure/watchdogs/media-pipeline-watchdog.timer
[Unit]
Description=Run Media Pipeline Watchdog every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
