# core/media-stack/homelab/infrastructure/watchdogs/stack-watchdog.timer
[Unit]
Description=Run stack-watchdog every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
